# sys-apps/corefreq/files/corefreqd.initd
#!/sbin/openrc-run

command=corefreqd
command_args="-q"
pidfile=

name="../corefreqd"
description="CoreFreq Daemon"